# infra/datapurity-backend.service
[Unit]
Description=DataPurity Backend FastAPI Service
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/datapurity/backend
ExecStart=/opt/datapurity/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 4
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=datapurity-backend

Environment=PYTHONUNBUFFERED=1
Environment=ENVIRONMENT=production

[Install]
WantedBy=multi-user.target
